# cron/webuser.cron
* * * * * ${WEBDIR}/html/passbolt_api/bin/cake EmailQueue.sender # Created by IaC
